# dist/debian/python/postrm
#!/bin/sh

set -e

case "$1" in
    purge|remove|failed-upgrade|abort-install|abort-upgrade|disappear)
        [ -n /opt/python-viri ] && rm -rf /opt/python-viri
        ;;
    upgrade)
        invoke-rc.d virid stop || true
        invoke-rc.d virid start || true
        ;;
    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
esac

#DEBHELPER#

exit 0

